# metrics/bin/metrics
#!/bin/bash
#
# /etc/init.d/@proj_name@
#
# Startup script for @proj_name@ App
# Ubuntu-ized
# 
# chkconfig: 2345 20 80
# description: Starts and stops App
# pidfile: /var/run/@proj_name@/@proj_name@.pid

. /lib/lsb/init-functions

#export JAVA_HOME=/usr/java/default
export APP_NAME=metrics
export APP_HOME=/usr/local/share/$APP_NAME
export APP_OWNR=$APP_NAME
LOG_DIR=/var/log/$APP_NAME
LOG_FILE=$LOG_DIR/$APP_NAME.log

DISTRO=`lsb_release -i|awk -F':' '{print $2}' |tr -d [:blank:]`
if [ $DISTRO == "Ubuntu" ]; then
  LOCK_FILE=/run/lock/$APP_NAME
  PID_DIR=/run/$APP_NAME
else
  LOCK_FILE=/var/lock/subsys/$APP_NAME
  PID_DIR=/var/run/$APP_NAME
fi

PID_FILE=$PID_DIR/$APP_NAME.pid
APP_EXEC=$APP_HOME/bin/$APP_NAME.sh

RETVAL=0

if [ `id -u` -ne 0 ]; then
    echo "The $APP_NAME init script can only be run as root"
    exit 1
fi

create_dir() {
  if [ ! -d  $1 ]; then
    mkdir -p $1
    chown $APP_NAME:$APP_NAME $1
  fi
}

start() {
  log_daemon_msg "Starting Nexgen Service" "$APP_NAME"
  create_dir $LOG_DIR
  create_dir $PID_DIR
  su $APP_OWNR -c "$APP_EXEC -p $PID_FILE " > $LOG_FILE 2>&1
  RETVAL=$?
  if [ $RETVAL -eq 0 ]; then
    log_success_msg "$APP_NAME started"
    touch $LOCK_FILE
    echo
  else
    log_failure_msg "$APP_NAME start failed"
    RETVAL=1; exit $RETVAL;
  fi
}

stop() {
  log_daemon_msg "Shutting down Nexgen Service" "$APP_NAME"
  killproc -p $PID_FILE $APP_NAME
  if [ $RETVAL -eq 0 ]; then
    log_success_msg "$APP_NAME shutdown"
    rm -f $LOCK_FILE
  else
    log_failure_msg "$APP_NAME shutdown"
    echo  
    RETVAL=1; exit $RETVAL
  fi
  RETVAL=0
}

restart() {
  stop
  start
}


case "$1" in
  start)
    start
    RETVAL=$?
  ;;
  stop)
    stop
    RETVAL=$?
  ;;
  restart)
    restart
    RETVAL=$?
  ;;
  status)
    status -p $PID_FILE $APP_NAME
    RETVAL=$?
  ;;
  *)
    echo "Usage: `basename $0` start|stop|restart|status"
    exit 1
esac

exit $RETVAL
